# backend/crontab
# Every 2 hours
0 */2 * * * cd /app && python manage.py scrapeapartments